# sd64/usr/lib/systemd/system/sdserver.socket
[Unit]
Description=sd Per-Connection Server

[Socket]
ListenStream=/tmp/sdsys/sdserver.socket
Accept=true

[Install]
WantedBy=sockets.target
