# sd64/usr/lib/systemd/system/sdserver@.service
[Unit]
Description=sd Per-Connection Server

[Service]
User=root
Group=sdusers
ExecStart=/usr/local/sdsys/bin/sd -n
StandardInput=socket

